# etc/update-motd.d/00-header
#!/bin/sh

echo ""
echo "     Your ASCII logo here!"
echo ""
